# MusicaEnArranquePC (Despacito)/dspasito.service
[Unit]

Description= Está escuchando despacito

[Service]

ExecStart=/usr/bin/dspsasito.sh

[Install]

WantedBy=multi-user.target

